# debian/vmango.upstart
description "Vmango web service"

start on runlevel [2345]
stop on runlevel [016]

script
    if [ -f /etc/default/vmango ]; then
        . /etc/default/vmango
    fi

    exec /usr/bin/vmango --config=/etc/vmango/vmango.conf
end script
